# Docker/RPi-Node/crontab
* * * * * root /usr/local/bin/python /farmalytics/scripts/dht.py >> /farmalytics/logs/dht.log
* * * * * root /usr/local/bin/python /farmalytics/scripts/outside.py >> /farmalytics/logs/outside.log
* * * * * root /usr/local/bin/python /farmalytics/scripts/moisture.py >> /farmalytics/logs/moisture.log
